# priv/deploy/hostctl.service
[Unit]
Description=Hostctl Phoenix Server
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=hostctl
Group=hostctl
WorkingDirectory=/opt/hostctl
EnvironmentFile=/etc/hostctl/env
ExecStart=/opt/hostctl/bin/server
ExecStop=/opt/hostctl/bin/hostctl stop
ExecStartPre=+/usr/bin/mkdir -p /run/sudo/ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hostctl

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/hostctl /var/log/hostctl /usr/local/src/hostctl /run/sudo /var/cache/apt /var/lib/apt /var/lib/dpkg
RuntimeDirectory=hostctl
RuntimeDirectoryMode=0700

[Install]
WantedBy=multi-user.target
